# bin/org.channing.octopus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.channing.octopus</string>
    <key>RunAtLoad</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/python3</string>
	  <string>/Users/channing/dotfiles/bin/check-octopus.py</string>
    </array>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>StandardOutPath</key>
    <string>/Users/channing/Library/Logs/octopus.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/channing/Library/Logs/octopus.log</string>
</dict>
</plist>

